# motibeam.service
[Unit]
Description=MotiBeam OS - Multi-Vertical Ambient Computing Platform
Documentation=https://github.com/quamie25/MotiBeam-OS
After=graphical.target

[Service]
Type=simple
User=motibeam
WorkingDirectory=/home/motibeam/MotiBeam-OS
Environment="DISPLAY=:0"
ExecStart=/bin/bash -c 'cd /home/motibeam/MotiBeam-OS && source motibeam-env/bin/activate && python3 motibeam_app.py'
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
